# nuc_display.service
[Unit]
Description=NUC Display Dashboard
After=network.target graphics.target alsa-restore.service

[Service]
Type=simple
User=@USER@
WorkingDirectory=@WORKDIR@
ExecStart=@WORKDIR@/build/nuc_display
Restart=always
RestartSec=5

# Grant access to hardware even in headless mode
Environment=XDG_RUNTIME_DIR=/run/user/@UID@
Environment=GDK_BACKEND=drm
Environment=EGL_PLATFORM=gbm
# ALSA configuration
Environment=ALSA_CARD=@ALSA_CARD@

# Security and Sandbox (Optional but good practice)
# We need to keep some privileges for DRM/Render access
SupplementaryGroups=video render audio plugdev input

[Install]
WantedBy=multi-user.target
